# deploy/systemd/bookrec-train-automation.timer.template
[Unit]
Description=Run BookRec automation once per day

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
